# debian/cron.d
#
# Regular cron jobs for the tiger package
#
0 * * * *      root    /usr/sbin/tigercron -q
